# images/personal/playbooks/files/btrfs-snapshot@.timer
[Unit]
Description=Create btrfs snapshot daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
