# images/personal/playbooks/files/btrfs-snapshot@.service
[Unit]
Description= BTRFS Snapshot for %I

[Service]
ExecStart=/usr/bin/btrfs-snapshot --instance %I
