# images/personal/playbooks/files/wasync.service
[Unit]
Description=Sync home directory with wasabu S3 storage

[Service]
ExecStart=/usr/bin/wasync
